# packages/cbr-api-gateway/cbr_api_gateway-0.7.0.tar.gz/cbr_api_gateway-0.7.0/cbr_api_gateway/nginx/cbr-api-gateway.supervisor.conf
[program:cbr_api_gateway]
#directory=%(ENV_SITE_PACKAGES_DIR)/cbr_api_gateway/nginx       ; this was not working, although SITE_PACKAGES_DIR=/usr/local/lib/python3.12/site-packages
#command=%(ENV_SITE_PACKAGES_DIR)/cbr_api_gateway/nginx/entrypoint.sh
#directory=/usr/local/lib/python3.12/site-packages/cbr_api_gateway/nginx
#command=/usr/local/lib/python3.12/site-packages/cbr_api_gateway/nginx/entrypoint.sh
command=/app/cbr_api_gateway/entrypoint.sh
directory=/app/cbr_api_gateway

autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0